# docker-images/kafka/supervisor/kafka-init.conf
[program:kafka-init]
command=/initialize-entrypoint.sh
autostart=true
autorestart=false